# voltage-websocket.service
[Unit]
Description=Volt WebSocket Server
After=network.target

[Service]
Type=simple
User=volt
WorkingDirectory=/home/bluet/Documents/EIVoltChat/voltchat/Voltage
ExecStart=/usr/bin/node /home/bluet/Documents/EIVoltChat/voltchat/Voltage/microservices/websocket.js
Restart=always
RestartSec=10
Environment=NODE_ENV=production
Environment=VOLT_SERVICE=websocket

[Install]
WantedBy=multi-user.target
